# infra/luckyniuma-backend.service
[Unit]
Description=LuckyNiuMaNote Express（站点与 API）
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/LuckyNiuMaNote
Environment=NODE_ENV=production
Environment=LISTEN_HOST=127.0.0.1
Environment=TRUST_PROXY=1
Environment=PORT=3000
Environment=SITE_PUBLIC_URL=https://luckyniuma.com
ExecStart=/usr/bin/node /home/ubuntu/LuckyNiuMaNote/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
